# cron/open-meteo.crontab
# Timezone for schedule fields below
TZ=Asia/Ho_Chi_Minh
APP_ENV=production

# STSPortal PHP app root (CodeIgniter — `php spark weather:open-meteo-pull`)
STSPORTAL_DIR=/var/www/STSPortal

# Logs + lock
LOG_DIR=/home/stsvps/stsportal/var/log/stsweather

# Daily 02:00 Asia/Ho_Chi_Minh — spark CLI, all active catalog locations (no HTTP / curl)
0 2 * * * /bin/bash -lc 'mkdir -p /home/stsvps/stsportal/var/run /home/stsvps/stsportal/var/log && flock -n /home/stsvps/stsportal/var/run/cron_open_meteo.lock /bin/bash /home/stsvps/stsportal/scripts/cron_open_meteo_scan.sh >> /home/stsvps/stsportal/var/log/cron_open_meteo_runner.log 2>&1'

# --- One-liner without wrapper script (optional) ---
# 0 2 * * * flock -n /home/stsvps/stsportal/var/run/cron_open_meteo.lock bash -lc 'cd /home/stsvps/STSPortal && php spark weather:open-meteo-pull >> /home/stsvps/stsportal/var/log/stsweather/spark_open_meteo_$(date +\%Y\%m\%d).log 2>&1'
